# ansible-experiment/java-app-experiment/deployment/myapp.service
[Unit]
Description={{app_name}}
After=syslog.target

[Service]
User=root
ExecStart=/var/{{app_name}}/{{app_name}}.war
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
